# .restic/restic.timer
[Unit]
Description=Backup with restic daily

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=900

[Install]
WantedBy=timers.target
